# files/confluence.conf
description "Use Upstart to kick off confluence auto-magic-ally :)"
author "Reich"
start on runlevel [2345]
exec /home/ubuntu/atlassian/confluence/bin/start-confluence.sh

